# als-training-service/crontab
0 2 * * 0 python /app/train.py >> /var/log/cron.log 2>&1
